# files/supervisord.conf
[supervisord]
nodaemon=true

[program:app]
command=npm run start 
autostart=true
autorestart=true


